# cron/crontab
#crontab
0 8 * * * root **applications/seacemas
